# ansible/deploy/roles/omniapp_second/templates/supervisord.conf
[unix_http_server]
file = /var/tmp/supervisor.sock ;

[supervisord]
;http_port=/var/tmp/supervisor.sock ; (default is to run a UNIX domain socket server)
;http_port=127.0.0.1:9001  ; (alternately, ip_address:port specifies AF_INET)
;sockchmod=0700              ; AF_UNIX socketmode (AF_INET ignore, default 0700)
;sockchown=nobody.nogroup     ; AF_UNIX socket uid.gid owner (AF_INET ignores)
;umask=022                   ; (process file creation umask;default 022)
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;http_username=user          ; (default is no username (open system))
;http_password=123           ; (default is no password (open system))
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;user=chrism                 ; (default is current user, required if root)
;directory=/tmp              ; (default is not to cd during start)
;environment=KEY=value       ; (key value pairs to add to environment)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
    
[supervisorctl]
serverurl=unix:///var/tmp/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")

; The below sample program section shows all possible program subsection values,
; create one or more 'real' program: sections to be able to control them under
; supervisor.

;[program:theprogramname]
;command=/bin/cat            ; the program (relative uses PATH, can take args)
;priority=999                ; the relative start priority (default 999)
;autostart=true              ; start at supervisord start (default: true)
;autorestart=true            ; retstart at unexpected quit (default: true)
;startsecs=10                ; number of secs prog must stay running (def. 10)
;startretries=3              ; max # of serial start failures (default 3)
;exitcodes=0,2               ; 'expected' exit codes for process (default 0,2)
;stopsignal=QUIT             ; signal used to kill process (default TERM)
;stopwaitsecs=10             ; max num secs to wait before SIGKILL (default 10)
;user=chrism                 ; setuid to this UNIX account to run the program
;log_stdout=true             ; if true, log program stdout (default true)
;log_stderr=true             ; if true, log program stderr (def false)
;logfile=/var/log/cat.log    ; child log path, use NONE for none; default AUTO
;logfile_maxbytes=1MB        ; max # logfile bytes b4 rotation (default 50MB)
;logfile_backups=10          ; # of logfile backups (default 10)

#~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~
#
# Celery Workers signals
# http://docs.celeryproject.org/en/latest/userguide/workers.html#process-signals
#
# TERM	Warm shutdown, wait for tasks to complete.
# QUIT	Cold shutdown, terminate ASAP
# USR1	Dump traceback for all active threads.
# USR2	Remote debug, see celery.contrib.rdb.
#
#~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~


[program:fts-celery-worker-esperar-finaliza-campana]
command = /bin/env
	PYTHONPATH={{ install_prefix }}ominicontacto:{{ install_prefix }}local
	DJANGO_SETTINGS_MODULE=fts_web.settings
	FTS_LOGFILE=celery-worker-esperar-finaliza-campana.log
	{{ install_prefix }}virtualenv/bin/celery
	-A fts_daemon.fts_celery_daemon.ominicontacto
	--no-color
	--logfile={{ install_prefix }}log/celery-esperar-finaliza-campana.log
	--loglevel=info
	--queues=esperar_y_finalizar_campana
	--concurrency=4
	--uid=ftsender
	--gid=ftsender
	--maxtasksperchild=1
	worker
autostart=true              ; start at supervisord start (default: true)
autorestart=true            ; retstart at unexpected quit (default: true)
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=180            ; max num secs to wait before SIGKILL (default 10)
priority=999
stdout_logfile={{ install_prefix }}log/celery-esperar-finaliza-campana.log
stderr_logfile={{ install_prefix }}log/celery-esperar-finaliza-campana.err

[program:fts-celery-worker-finalizar-campana]
command = /bin/env
	PYTHONPATH={{ install_prefix }}ominicontacto:{{ install_prefix }}local
	DJANGO_SETTINGS_MODULE=fts_web.settings
	FTS_LOGFILE=celery-worker-finalizar-campana.log
	{{ install_prefix }}virtualenv/bin/celery
	-A fts_daemon.fts_celery_daemon.ominicontacto
	--no-color
	--logfile={{ install_prefix }}log/celery-finalizar-campana.log
	--loglevel=info
	--queues=finalizar_campana
	--concurrency=1
	--uid=ftsender
	--gid=ftsender
	--maxtasksperchild=1
	worker
autostart=true              ; start at supervisord start (default: true)
autorestart=true            ; retstart at unexpected quit (default: true)
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
stopsignal=TERM             ; signal used to kill process (default TERM)
stopwaitsecs=180            ; max num secs to wait before SIGKILL - 180 = 3 minutos
priority=999
stdout_logfile={{ install_prefix }}log/celery-finalizar-campana.log
stderr_logfile={{ install_prefix }}log/celery-finalizar-campana.err

#~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~
#
# Daemons
#
#~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~~


[program:fts-llamador-poll-daemon]
command = /bin/env
	PYTHONPATH={{ install_prefix }}ominicontacto:{{ install_prefix }}local
	DJANGO_SETTINGS_MODULE=fts_web.settings
	FTS_LOGFILE=daemon-llamador-poll.log
	{{ install_prefix }}virtualenv/bin/python
	{{ install_prefix }}ominicontacto/fts_daemon/poll_daemon/main.py
autostart=true              ; start at supervisord start (default: true)
autorestart=true            ; retstart at unexpected quit (default: true)
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
stopsignal=TERM             ; signal used to kill process (default TERM)
user=ftsender
# stopwaitsecs=180            ; max num secs to wait before SIGKILL (default 10)
priority=500
stdout_logfile={{ install_prefix }}log/daemon-llamador-poll.log
stderr_logfile={{ install_prefix }}log/daemon-llamador-poll.err

[program:fts-fastagi-daemon]
command = /bin/env
	PYTHONPATH={{ install_prefix }}ominicontacto:{{ install_prefix }}local
	DJANGO_SETTINGS_MODULE=fts_web.settings
	FTS_LOGFILE=daemon-fastagi.log
	{{ install_prefix }}virtualenv/bin/python
	{{ install_prefix }}ominicontacto/fts_daemon/fastagi_daemon.py
autostart=true              ; start at supervisord start (default: true)
autorestart=true            ; retstart at unexpected quit (default: true)
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
stopsignal=TERM             ; signal used to kill process (default TERM)
user=ftsender
# stopwaitsecs=180            ; max num secs to wait before SIGKILL (default 10)
priority=0
stdout_logfile={{ install_prefix }}log/daemon-fastagi.log
stderr_logfile={{ install_prefix }}log/daemon-fastagi.err

[program:fts-chequeador-campanas-vencidas]
command = /bin/env
	PYTHONPATH={{ install_prefix }}ominicontacto:{{ install_prefix }}local
	DJANGO_SETTINGS_MODULE=fts_web.settings
	FTS_LOGFILE=daemon-chequeador-campanas-vencidas.log
	{{ install_prefix }}virtualenv/bin/python
	{{ install_prefix }}ominicontacto/fts_daemon/services/finalizador_vencidas_daemon.py
autostart=true              ; start at supervisord start (default: true)
autorestart=true            ; retstart at unexpected quit (default: true)
startsecs=10                ; number of secs prog must stay running (def. 10)
startretries=3              ; max # of serial start failures (default 3)
stopsignal=TERM             ; signal used to kill process (default TERM)
user=ftsender
# stopwaitsecs=180            ; max num secs to wait before SIGKILL (default 10)
priority=999
stdout_logfile={{ install_prefix }}log/daemon-chequeador-campanas-vencidas.log
stderr_logfile={{ install_prefix }}log/daemon-chequeador-campanas-vencidas.err